# systemd/ccproxy.service
[Unit]
Description=ccProxy - OpenAI/Anthropic to Codex Responses proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# FHS layout: binary in /usr/local/bin, config in /etc/ccproxy, data in /var/lib/ccproxy.
ExecStart=/usr/local/bin/ccproxy
WorkingDirectory=/var/lib/ccproxy
EnvironmentFile=/etc/ccproxy/.env
# auth.json + SQLite db live under the data dir (HOME/.codex, HOME/.ccproxy).
Environment=HOME=/var/lib/ccproxy
Restart=on-failure
RestartSec=2
DynamicUser=no

[Install]
WantedBy=multi-user.target
